# rig-common-systemd/ai-rig-backup.timer
[Unit]
Description=AI Rig - Timer backup giornaliero

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
