# rig-common/systemd/ai-rig-backup.service
[Unit]
Description=AI Rig - Backup periodico config su 4° disco
[Service]
Type=oneshot
ExecStart=/usr/local/bin/70-backup.sh
